# scripts/cron/rcas6056-starreco.cron
45  05,08,11,14,17,20,23,02  * * * /star/u/starreco/bin/Auto-Spider.csh Check
