# scripts/cron/rcas6003-starlib.cron
* * * * * /star/u/starlib/bin/token.csh
05,35 * * * *   /star/u/starlib/bin/PDSpace.csh /afs/rhic.bnl.gov/star/doc/www/html/tmp/pub/overall.html

00 21 * * 2,4,6 /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh Clean
05 23 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh
00 04 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh icc
00 06 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh Insure

01,16,31,46 * * * * /afs/rhic.bnl.gov/star/packages/scripts/cvslog2html.pl
45 00,12,17   * * * /afs/rhic.bnl.gov/star/packages/scripts/swguide-cron.pl .dev dev
00 03         * * * /afs/rhic.bnl.gov/star/packages/scripts/swguide-cron.pl new pro old
00 04         * * * /afs/rhic.bnl.gov/star/packages/scripts/doxycron.pl >&/dev/null

# Reference for AddThisNode.csh and nova daemon added on Wed Nov 24 22:49:22 EST 2004
51 22,0,2,4,6,8,10,12,14,16,18,20 * * * /tmp/nova/novaStart.csh
43 0,2,4,6,8,10,12,14,16,18,20,22 * * * /usr/bin/test -e /star/u/starlib/novabin/AddThisNode.csh && /star/u/starlib/novabin/AddThisNode.csh 
