# scripts/cron/rcas6003-starlib.cron
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.22688 installed on Wed Jan 28 17:52:17 2004)
# (Cron version -- $Id: rcas6003-starlib.cron,v 1.5 2004/03/21 16:14:09 jeromel Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/afs/rhic/star/users/starlib/novabin/cron/rcas6003.rcf.bnl.gov.cron installed on Tue Oct 28 10:35:47 2003)
# (Cron version -- $Id: rcas6003-starlib.cron,v 1.5 2004/03/21 16:14:09 jeromel Exp $)

05,35 * * * *   /afs/rhic.bnl.gov/star/packages/scripts/PDSpace.pl /afs/rhic.bnl.gov/star/doc/www/html/tmp/pub/overall.html

00 21 * * 2,4,6 /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh Clean
00 06 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh Insure
05 23 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh
*/10 * * * *     /afs/rhic.bnl.gov/star/packages/scripts/TokenTest.csh

01,16,31,46 * * * * /afs/rhic.bnl.gov/star/packages/scripts/cvslog2html.pl
45 00,12,17   * * * /afs/rhic.bnl.gov/star/packages/scripts/swguide-cron.pl .dev dev
00 03         * * * /afs/rhic.bnl.gov/star/packages/scripts/swguide-cron.pl new pro old
00 04         * * * /afs/rhic.bnl.gov/star/packages/scripts/doxycron.pl >&/dev/null

37 10,12,14,16,18,20,22,0,2,4,6,8 * * * /tmp/nova/novaStart.csh
