# scripts/cron/rcas6041-starreco.cron
25  02,05,08,11,14,17,20,23  * * * /star/u/starreco/bin/Auto-Spider.csh Scan
