# scripts/cron/rmine606-starreco.cron
20 * * * * /star/u/starreco/bin/Auto-Spider
