# scripts/cron/rmine609-starlib.cron
# Reference for AddThisNode.csh and nova daemon added on Thu Nov 25 20:50:53 EST 2004
52 20,22,0,2,4,6,8,10,12,14,16,18 * * * /tmp/nova/novaStart.csh
55 0,2,4,6,8,10,12,14,16,18,20,22 * * * /usr/bin/test -e /star/u/starlib/novabin/AddThisNode.csh && /star/u/starlib/novabin/AddThisNode.csh 
