# scripts/cron/rmine603-starlib.cron
# This entry is coupled with rmine601's entries
#7,37         *   * * * /afs/rhic/star/packages/scripts/showquotas -html -o /afs/rhic/star/doc/www/html/tmp/pub/data03.html /star/data03
9,39         *   * * * /afs/rhic/star/packages/scripts/showquotas -html -o /afs/rhic/star/doc/www/html/tmp/pub/data05.html /star/data05

