# scripts/cron/rmine612-starreco.cron
26 * * * * /star/u/starreco/bin/Auto-Spider
