# scripts/cron/rplay27-starlib.cron
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.17196 installed on Mon Jun  7 12:29:05 2004)
# (Cron version -- $Id: rplay27-starlib.cron,v 1.2 2005/11/08 15:30:22 jeromel Exp $)
00 05 * * *     /afs/rhic.bnl.gov/star/packages/scripts/AutoBuild.csh Linux72

